# .config/systemd/user/media-x-rclone-aar.mount
[Mount]
What=aar:
Where=/media/x/rclone/aar

[Install]
WantedBy=default.target
