# .config/systemd/user/media-nv-rclone-aar.mount
[Mount]
What=aar:
Where=/media/nv/rclone/aar

[Install]
WantedBy=default.target