# .config/systemd/user/spotlight.timer
[Unit]
Description=Runs Spotlight daily

[Timer]
OnBootSec=1m
OnUnitActiveSec=2h
Unit=spotlight.service

[Install]
WantedBy=timers.target