# .config/systemd/user/spotlight.timer
[Unit]
Description=Runs Spotlight daily

[Timer]
OnCalendar=daily
Persistent=yes
Unit=spotlight.service

[Install]
WantedBy=timers.target
