# .config/systemd/user/mpd.service.d/override.conf
[Unit]
Before=mpDris2.service
BindsTo=mpDris2.service

[Install]
Also=mpd.socket