# .config/systemd/user/mpd.service.d/override.conf
[Unit]
Before=mpdris2-rs.service
BindsTo=mpdris2-rs.service

[Install]
Also=mpd.socket
