# .config/autostart/ssh-add.desktop
# ~/.config/autostart/ssh-add.desktop
[Desktop Entry]
Exec=ssh-add .ssh/keys/aur
Name=ssh-add
Type=Application
Terminal=false
