# .config/autostart/warp-taskbar.desktop
[Desktop Entry]
Exec=/usr/bin/warp-taskbar
Icon=dialog-scripts
Name=warp-taskbar
Type=Application
X-KDE-AutostartScript=true
