# .config/systemd/user/app-discord@autostart.service.d/override.conf
[Service]
ExecStart=
ExecStart=discord --start-minimized --enable-features=VaapiVideoDecoder --disable-features=UseChromeOSDirectVideoDecoder --ignore-gpu-blacklist
Restart=on-failure